# hatespeech_core/data/search_streaming/supervisord.conf
; Sample supervisor config file for daemonizing the twitter search to mongodb software
; You must generate your own consumer and token keys and secrets
; The one included is only for illustration purposes

[inet_http_server]          ; inet (TCP) server disabled by default
port=127.0.0.1:9001         ; (ip_address:port specifier, *:port for all iface)
username=manorastroman      ; (default is no username (open server))
password=kingofthedragonmen ; (default is no password (open server))

[supervisord]
stopsignal=INT
logfile=supervisord.log      ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=supervisord.pid      ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
username=manorastroman       ; should be same as http_username if set
password=kingofthedragonmen  ; should be same as http_password if set

[program:tapsearch]
command=tap search --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 -q "@realDonaldTrump OR @PreppyAltRight OR @SomeSchmick OR @pulpfascist OR @MrSalty_2" -v DEBUG
stdout_logfile=tap_search.log
stderr_logfile=tap_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapsearch]
command=tap search --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 -q "@nostariun OR @CantStumpTrump1 OR @AnthonyGillis OR @pulpfascist OR @_AltRight_Anew" -v DEBUG
stdout_logfile=tap_search.log
stderr_logfile=tap_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapsearch]
command=tap search --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 -q "@Vendetta92429 OR @stay___woke OR @Night_0f_Fire OR @AdlerHess_zwei OR @White_Jesus666" -v DEBUG
stdout_logfile=tap_search.log
stderr_logfile=tap_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapsearch]
command=tap search --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 -q "@r0ma_invicta OR @BlaineBruce OR @CarlKenner OR @SirAnglo OR @A1yosha" -v DEBUG
stdout_logfile=tap_search.log
stderr_logfile=tap_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapsearch]
command=tap search --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 -q "@1OfTheGoodGoys OR @Dr_Whomsoever OR @musicasacra62 OR @altRightCH" -v DEBUG
stdout_logfile=tap_search.log
stderr_logfile=tap_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapstream]
command=tap stream --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 --access-token 232890729-B8emFPeCTngXPSeJSSlP3FaZdPnCGGx9VoJwla73 --access-token-secret ZNvg11PKVLROQ6qHG5y8j1qiV9AGVJ7J0KCXnwN82vfS3 --track "'deplorable','republican', 'nationalist', 'alt right', 'anglo', 'whitegenocide', '#blacklivesmatter', '#bluelivesmatter', '#alllivesmatter','white power', 'obomber', 'democrat', 'great again', 'dirty mexicans', 'muzzie', 'Jihadi', 'Zionazi', 'khazar', 'camel fucker', 'hymie', 'papist', 'heeb', 'Christ killer', 'Hun', 'hebro', 'higger', 'soup taker ', 'hebe', 'oven dodger', 'proddy dog', 'African catfish', 'dogun', 'Amo', 'mackerel snapper ', 'carrot snapper', 'four by two', 'proddywoddy', 'Russellite ', 'Bengali', 'clam', 'clamhead', 'dogan', 'Gerudo', 'gew', 'Orangie', 'prod', 'proddywhoddy', 'snout', 'spike', 'WASP', 'WIC', 'wigger', 'wigga', 'surrender monkey', 'Zionazi', 'yardie', 'whigger', 'pommy', 'pommie', 'sideways vagina', 'wog', 'polack', 'wexican', 'wop', 'sideways pussy', 'shanty Irish', 'white nigger', 'soup taker ', 'roofucker', 'smoked Irish', 'wiggerette', 'sideways cooter', 'whore from Fife', 'pom', 'pommie grant', 'Punjab', 'Rhine monkey', 'sawney', 'seppo', 'septic', 'shant', 'sheepfucker', 'Shelta', 'skinny', 'Skip', 'Skippy', 'smoked Irishman', 'snout', 'sole', 'spink', 'squarehead', 'Taffy', 'tan', 'Tommy', 'tunnel digger', 'Yank', 'Yankee', 'cunt', 'twat', 'tranny', 'dyke', 'bint', 'bitch', 'pussy', 'bird', 'ho', 'hoe', 'scag', 'skag', 'cotton picker', 'wigger', 'hillbilly', 'wigga', 'surrender monkey', 'Chinaman', 'goober', 'whigger', 'spiv', 'zipperhead', 'kraut', 'paleface', 'Hun', 'white nigger', 'spook', 'smoked Irish', 'wiggerette', 'sooty', 'boojie', 'buffie', 'celestial', 'Charlie', 'cheese eating surrender monkey', 'crow', 'dink', 'div', 'divvy', 'fuzzy wuzzy', 'sawney', 'sheeny', 'shine', 'shiner', 'smoked Irishman', 'snowflake', 'squarehead', 'teapot', 'Tommy', 'tunnel digger', 'WIC', 'zip', 'zippohead', 'retard', 'cripple', 'sperg', 'retarded', 'idiot', 'mong', 'knacker', 'greaseball', 'pikey', 'hoser', 'Chinaman', 'bogan', 'plastic paddy', 'guala', 'limey', 'camel fucker', 'ocker', 'Argie', 'piker', 'paki', 'chink', 'chonky', 'kraut', 'Jap', 'booner', 'paddy', 'bog trotter', 'Hun', 'banana bender', 'eyetie', 'blaxican', 'Bog Irish', 'mick', 'ginzo', 'dingo fucker', 'Lebbo', 'Leb', 'Gwat', 'japie', 'island nigger', 'lugan', 'ice nigger', 'pogue', 'lowlander', 'guala guala', 'eh hole', 'guido', 'dogun', 'Hunyak', 'bog hopper', 'pohm', 'Honyock', 'cow kisser', 'ice monkey', 'greaser', 'bamboo coon', 'bog jumper', 'buk buk', 'dyke jumper', 'fresh off the boat', 'guinea', 'Hunyock', 'buckethead', 'Chinese wetback', 'chinig', 'leprechaun', 'ABC', 'banana lander', 'bans and cans', 'beaner shnitzel', 'Bengali', 'bhrempti', 'black dago', 'blockhead', 'book book', 'boon', 'boxhead', 'bubble', 'butterhead', 'celestial', 'Charlie', 'cheese eating surrender monkey', 'chigger', 'chink a billy', 'chunky', 'dink', 'dogan', 'domes', 'dot head', 'FOB', 'fog nigger', 'gew', 'heinie', 'Honyak', 'Hunkie', 'Hunky', 'Jerry', 'mack', 'Merkin', 'mickey', 'Mickey Finn', 'mung', 'nip', 'Pepper', 'Pepsi', 'piky', 'white trash', 'trailer trash', 'trailer park trash', 'hoser', 'bludger', 'hayseed', 'bogan', 'yardie', 'yobbo', 'dole bludger', 'charver', 'redneck', 'cracker', 'spiv', 'chav', 'charva', 'skanger', 'hick', 'scally', 'steek', 'gurrier', 'jhant', 'scanger', 'jant', 'booner', 'bog trotter', 'spide', 'Bog Irish', 'scallie', 'moke', 'millie', 'hoosier', 'stump jumper', 'senga', 'bog hopper', 'carrot snapper', 'mil bag', 'bog jumper', 'cave nigger', 'conspiracy theorist', 'ned', 'bitter clinger', 'boojie', 'boon', 'chink a billy', 'div', 'divvy', 'rube', 'scag', 'shant', 'shit kicker', 'skag', 'slag', 'trash', 'tree jumper', 'yob', 'faggot', 'fag', 'tranny', 'dyke', 'queer', 'fruit', 'Fairy', 'queen'" -v DEBUG
stdout_logfile=tap_stream.log
stderr_logfile=tap_stream_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapstream]
command=tap stream --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 --access-token 232890729-B8emFPeCTngXPSeJSSlP3FaZdPnCGGx9VoJwla73 --access-token-secret ZNvg11PKVLROQ6qHG5y8j1qiV9AGVJ7J0KCXnwN82vfS3 --track "'knacker', 'niggur', 'nigger', 'peckerwood', 'nigga', 'nig nog', 'pikey', 'nicca', 'niggar', 'nigglet', 'nigguh', 'niggah', 'niglet', 'moon cricket', 'plastic paddy', 'mud duck', 'Kushite', 'niggress', 'nig', 'mud shark', 'munter', 'khazar', 'pickaninny', 'ofay', 'piker', 'paki', 'octaroon', 'moulie', 'kyke', 'paleface', 'niggor', 'nigette', 'Oriental', 'lemonhead', 'mulignan', 'paddy', 'nigra', 'mangia cake', 'kike', 'mulato', 'mockey', 'moulinyan', 'ling ling', 'latrino', 'Kushi', 'mick', 'moulignon', 'mud person', 'nigar', 'nigre', 'nigor', 'niger', 'muk', 'ping pang', 'moke', 'oven dodger', 'pogue', 'moky', 'mocky', 'mosshead', 'mockie', 'pineapple nigger', 'octroon', 'pancake face', 'kotiya', 'lawn jockey', 'moss eater', 'mutt', 'negro', 'nitchee', 'Junior Mint', 'lefty', 'mack', 'mickey', 'Mickey Finn', 'moch', 'mock', 'monkey', 'Moor', 'moxy', 'muktuk', 'mung', 'munt', 'nacho', 'neche', 'neechee', 'neejee', 'net head', 'nichi', 'nichiwa', 'nidge', 'nitchie', 'nitchy', 'Oreo', 'pancake', 'papoose', 'Pepper', 'Pepsi', 'piky', 'pinto', 'yellow bone', 'white trash', 'wigger', 'wigga', 'yokel', 'wetback', 'yardie', 'whigger', 'whitey', 'zipperhead', 'ZOG', 'wog', 'wexican', 'white nigger', 'zigabo', 'wiggerette', 'ZOG lover', 'WASP', 'white chocolate', 'WIC', 'wink', 'yellow', 'yid', 'zebra', 'zip', 'zippohead', 'uncivilised', 'raghead', 'shyster', 'Uncle Tom', 'porch monkey', 'trailer trash', 'trailer park trash', 'race traitor', 'quadroon', 'red bone', 'tar baby', 'spear chucker', 'Rico Suave', 'sheister', 'Shylock', 'sideways vagina', 'slant eye', 'shit heel', 'spic', 'towel head', 'popolo', 'squaw', 'touch of the tar brush', 'redskin', 'sideways pussy', 'wagon burner', 'shanty Irish', 'roundeye', 'redlegs', 'sambo', 'prairie nigger', 'spigger', 'thicklips', 'slopehead', 'quashie', 'roofucker', 'timber nigger', 'poppadom', 'spigotty', 'semihole', 'spook', 'uncircumcised baboon', 'smoked Irish', 'spice nigger', 'Velcro head', 'tinkard', 'uncivilized', 'sideways cooter', 'sooty', 'spick', 'spickaboo', 'pointy head', 'pollo', 'powderburn', 'property', 'Punjab', 'shade', 'sheeny', 'Shelta', 'shine', 'shiner', 'Shy', 'Skip', 'Skippy', 'slant', 'slit', 'slope', 'slopey', 'slopy', 'smoke jumper', 'smoked Irishman', 'snowflake', 'sole', 'spade', 'spig', 'spik', 'spink', 'squarehead', 'squinty', 'stovepipe', 'sub human', 'sucker fish', 'suntan', 'Taffy', 'teapot', 'tenker', 'tiger', 'tincker', 'tinkar', 'tinker', 'tinkere', 'trash', 'tunnel digger', 'Twinkie', 'tyncar', 'tynekere', 'tynkard', 'tynkare', 'tynker', 'tynkere', 'coon ass', 'bogan', 'anchor baby', 'camel jockey', 'ching chong', 'alligator bait', 'border jumper', 'beaner', 'border hopper', 'camel fucker', 'coon', 'bix nood', 'af', 'chonky', 'coolie', 'Christ killer', 'camel humper', 'Buckwheat', 'boonie', 'bluegum', 'blaxican', 'border bunny', 'Africoon', 'buckra', 'abo', 'Americoon', 'boong', 'Afro-Saxon', 'banjo lips', 'burrhead', 'carpet pilot', 'African catfish', 'Caublasian', 'camel cowboy', 'Armo', 'bean dipper', 'border nigger', 'chili shitter', 'beach nigger', 'cab nigger', 'Buddhahead', 'bamboo coon', 'brass ankle', 'cave nigger', 'buk buk', 'camel jacker', 'albino', 'ape', 'Aunt Jemima', 'boo', 'buckethead', 'Chinese wetback', 'chinig', 'chug', 'abbo', 'ABC', 'ABCD', 'Angie', 'Anglo', 'Ann', 'apple', 'Aunt Jane', 'Aunt Mary', 'Aunt Sally', 'azn', 'banana', 'bans and cans', 'beaner shnitzel', 'beaney', 'Bengali', 'bhrempti', 'black Barbie', 'black dago', 'bong', 'boojie', 'book book', 'boonga', 'Bounty bar', 'brownie', 'bubble', 'buck', 'buffie', 'bug eater', 'bumblebee', 'bung', 'bunga', 'butterhead', 'can eater', 'celestial', 'Charlie', 'chee chee', 'chi chi', 'chigger', 'chink a billy', 'chunky', 'cocoa', 'Cocoa Puff', 'coconut', 'colored', 'coloured', 'gypo', 'greaseball', 'Jihadi', 'cotton picker', 'darkie', 'honkey', 'jiggaboo', 'honkie', 'curry muncher', 'house nigger', 'eurotrash', 'honky', 'half breed', 'guala', 'gyppo', 'jungle bunny', 'jigaboo', 'gator bait', 'gippo', 'Cushite', 'cracker', 'dune coon', 'gyp', 'darky', 'gook', 'hymie', 'groid', 'jigger', 'injun', 'cowboy killer', 'heeb', 'jigga', 'diaper head', 'Jap', 'jigg', 'gubba', 'gooky', 'darkey', 'hebro', 'jiggabo', 'halfrican', 'Cushi', 'eyetie', 'higger', 'Jewbacca', 'gub', 'dinge', 'ginzo', 'japie', 'island nigger', 'hebe', 'ice nigger', 'dago', 'hairyback', 'guala guala', 'guido', 'gyppy', 'dune nigger', 'ghetto hamster', 'cow kisser', 'curry slurper', 'four by two', 'greaser', 'gin jockey', 'gun burglar', 'golliwog', 'gook eye', 'guinea', 'Jim Fish', 'curry stinker', 'eight ball', 'ginger', 'ikey mo', 'jigarooni', 'crow', 'dego', 'dink', 'domes', 'dot head', 'egg', 'eggplant', 'fez', 'fog nigger', 'frog', 'fuzzy', 'fuzzy wuzzy', 'gable', 'gew', 'ghetto', 'ghost', 'gin', 'gipp', 'gyppie', 'Hunkie', 'Hunky', 'ike', 'ikey', 'iky', 'jig', 'jijjiboo'" -v DEBUG
stdout_logfile=tap_stream.log
stderr_logfile=tap_stream_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

[program:tapstream]
command=tap stream --consumer-key 82ouMNTDovEVsUJAKaayoPVZ9 --consumer-secret TSZ86pyseQd1PpRbv5uH1trYYY6qbq0nHZLBbaMPXzehUDf2a6 --access-token 232890729-B8emFPeCTngXPSeJSSlP3FaZdPnCGGx9VoJwla73 --access-token-secret ZNvg11PKVLROQ6qHG5y8j1qiV9AGVJ7J0KCXnwN82vfS3 --track "realDonaldTrump,PreppyAltRight,SomeSchmick,pulpfascist,MrSalty_2,nostariun,CantStumpTrump1,AnthonyGillis,pulpfascist,_AltRight_Anew,Vendetta92429,stay___woke,Night_0f_Fire,AdlerHess_zwei,White_Jesus666,r0ma_invicta,BlaineBruce,CarlKenner,SirAnglo,A1yosha,1OfTheGoodGoys,Dr_Whomsoever,musicasacra62,altRightCH" -v DEBUG
stdout_logfile=tap_stream.log
stderr_logfile=tap_stream_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10

